# deploy/dc-gra-vt-dashboard.service
[Unit]
Description=DC-Gra-vt-bot Dashboard API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/DC
EnvironmentFile=/home/ubuntu/DC/.env
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/usr/bin/test -f /home/ubuntu/DC/.env
ExecStartPre=/usr/bin/test -x /home/ubuntu/DC/.venv/bin/uvicorn
ExecStart=/home/ubuntu/DC/.venv/bin/uvicorn dashboard_api:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=8
User=ubuntu

[Install]
WantedBy=multi-user.target
