# deploy/dc-gra-vt-bot.service
[Unit]
Description=DC-Gra-vt-bot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dc-gra-vt-bot
EnvironmentFile=/opt/dc-gra-vt-bot/.env
ExecStart=/opt/dc-gra-vt-bot/.venv/bin/python bot.py
Restart=always
RestartSec=8
User=ubuntu

[Install]
WantedBy=multi-user.target
